# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# MySQL Database Service
[program:mysql]
command=/app/init-mysql.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/mysql.err.log
stdout_logfile=/var/log/supervisor/mysql.out.log
priority=100

# Nginx Web Server
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/nginx.err.log
stdout_logfile=/var/log/supervisor/nginx.out.log
priority=200

# Eureka Service Discovery
[program:eureka]
command=/app/start-eureka.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/eureka.err.log
stdout_logfile=/var/log/supervisor/eureka.out.log
priority=300
startsecs=30

# User Service
[program:user-service]
command=/app/start-user-service.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/user-service.err.log
stdout_logfile=/var/log/supervisor/user-service.out.log
priority=400
startsecs=45

# Content Service
[program:content-service]
command=/app/start-content-service.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/content-service.err.log
stdout_logfile=/var/log/supervisor/content-service.out.log
priority=500
startsecs=45